# get-iplayer.conf
[program:get-iplayer]
user = nobody
autorestart = False
command = /usr/bin/get-iplayer-script.sh