# get-iplayer.conf
[program:cronie]
user = root
command = /usr/bin/crond -n