# get-iplayer.conf
[program:get-iplayer]
user = nobody
startretries = 0
command = /usr/bin/get-iplayer-script.sh