# get-iplayer.conf
[program:cronie]
user = nobody
command = /usr/bin/crond -n