# get-iplayer.conf
[program:get-iplayer]
user = nobody
startsecs = 0
autorestart = False
command = /usr/bin/perl /usr/bin/get_iplayer.cgi --port=1935 --getiplayer=/usr/bin/get_iplayer
