# get-iplayer.conf
[program:get-iplayer]
user = nobody
startsecs = 0
autorestart = False
command = /usr/bin/get-iplayer-script.sh