# get_iplayer.conf
[program:get_iplayer-zack]
user = nobody
command = /usr/bin/get_iplayer-script.sh