# setup/get-iplayer.conf
[program:get-iplayer]
user = nobody
startsecs = 0
autorestart = False
command = /home/nobody/start.sh
umask = 000
